# debian/orchestrator/DEBIAN/postinst
#!/bin/sh

systemctl daemon-reload
systemctl enable arduino-orchestrator
systemctl start arduino-orchestrator
