# debian/orchestrator/DEBIAN/postinst
#!/bin/sh

chown -R arduino:arduino /home/arduino/.arduino-app-cli/examples

systemctl enable arduino-orchestrator
