# debian/orchestrator/DEBIAN/prerm
#!/bin/sh

systemctl disable arduino-orchestrator
systemctl stop arduino-orchestrator
